# ai-agents/portcl/init.d/portcl-agent
#!/sbin/openrc-run

name="portcl-agent"
description="PortCL AI Portage Optimization Agent"
command="/usr/local/bin/portcl-agent"
command_args="--config /etc/portcl/config.toml"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net
    after portage
    use logger
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /var/log/portcl
    checkpath --directory --owner root:root --mode 0755 /var/lib/portcl
}

start() {
    ebegin "Starting ${name}"
    start-stop-daemon --start --exec "${command}" \
        --pidfile "${pidfile}" \
        --background --make-pidfile \
        -- ${command_args}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop --pidfile "${pidfile}"
    eend $?
}